# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

/usr/sbin/nginx
